# scripts/systemd-timer/container-updates.timer
[Unit]
Description=Send container updates information

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

